# src/etc/cron.d/szen.pac
*/15 * * * * www-data /var/www/szen.pac.git/src/bin/update > /dev/null 2>&1

# vim: se ft=crontab fenc=utf-8 ff=unix:
